# swarm/nous-swarm-sync.timer
[Unit]
Description=NOUS Swarm sync med peers

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
