# scripts/nous-scraper.service
[Unit]
Description=NOUS Night Scraper — ugentlig søgning og ingest af nye dokumenter
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
User=nous
WorkingDirectory=/srv/nous/scripts
Environment=PATH=/srv/nous/app/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/nous
ExecStart=/srv/nous/app/.venv/bin/python3 /srv/nous/scripts/night_scraper.py
StandardOutput=append:/mnt/nous-data/logs/scraper.log
StandardError=append:/mnt/nous-data/logs/scraper.log
TimeoutStartSec=3600
